# scripts/dash-of-pi.service
[Unit]
Description=Dash of Pi - Raspberry Pi Dash Camera Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
StartLimitBurst=5

# User and group
User=dash-of-pi
Group=dash-of-pi

# Working directory
WorkingDirectory=/var/lib/dash-of-pi

# Executable
ExecStart=/usr/local/bin/dash-of-pi -config /etc/dash-of-pi/config.json -v

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dash-of-pi

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Security
NoNewPrivileges=true
PrivateTmp=true

# Device access for camera
PrivateDevices=false

[Install]
WantedBy=multi-user.target
