# scripts/dash-of-pi.service
[Unit]
Description=Dash of Pi - Raspberry Pi Dash Camera Service
After=network.target

[Service]
Type=simple

# Restart policy - automatically restart if the service crashes
Restart=on-failure
RestartSec=10

# Run as dedicated service user (not root) for security
User=dash-of-pi
Group=dash-of-pi

# Working directory for the service
WorkingDirectory=/var/lib/dash-of-pi

# Command to run
ExecStart=/usr/local/bin/dash-of-pi -config /etc/dash-of-pi/config.json

# Send logs to systemd journal (view with: journalctl -u dash-of-pi -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dash-of-pi

# Resource limits for Pi Zero 2W (512MB RAM)
# Max open files (videos, camera, HTTP connections)
LimitNOFILE=1024
# Max processes/threads
LimitNPROC=128
# Prevent OOM - limit to 200MB max (leaves room for system)
MemoryMax=200M
# Start throttling at 150MB
MemoryHigh=150M

# Security hardening
# Prevent privilege escalation
NoNewPrivileges=true
# Isolated /tmp directory
PrivateTmp=true

# Allow camera access (required for /dev/video0)
PrivateDevices=false

[Install]
WantedBy=multi-user.target
